# assets/zerod.service.in
# zerod user service — generated by `zerod service install`.
# To override any field without editing this file, drop a snippet at
# ~/.config/systemd/user/zerod.service.d/override.conf

[Unit]
Description=zerod headless audio / Bluetooth / systemd control daemon
Documentation=https://tangled.org/tsiry-sandratraina.com/zerod
After=network-online.target pipewire.service wireplumber.service
Wants=network-online.target

[Service]
Type=simple
ExecStart={{exec_start}}
Restart=on-failure
RestartSec=5

# Pin the bearer token via env so the random fallback isn't used at startup.
# Either set it inline here:
#   Environment=ZEROD_BEARER_TOKEN=changeme
# or — better — in a drop-in that's not regenerated by `service install`:
#   ~/.config/systemd/user/zerod.service.d/override.conf
#
# Configuration files (bind address, systemd unit allowlist, managed
# config registry) live in ~/.config/zerod/zerod.toml — see
# `zerod.toml.example` for the schema.

[Install]
WantedBy=default.target
