# assets/zerod.service.in
# zerod user service — generated by `zerod service install`.
# To override any field without editing this file, drop a snippet at
# ~/.config/systemd/user/zerod.service.d/override.conf

[Unit]
Description=zerod headless audio / Bluetooth / systemd control daemon
Documentation=https://tangled.org/tsiry-sandratraina.com/zerod
After=network-online.target pipewire.service wireplumber.service
Wants=network-online.target

[Service]
Type=simple
ExecStart={{exec_start}}
Restart=on-failure
RestartSec=5

# Bearer token pinned at install time so zerod doesn't generate a fresh
# random one on every restart (which would break any client that's
# already paired with the old value). To rotate it, re-run
# `zerod service install --force` (or `--token <value>` to pin a chosen
# value), then `systemctl --user restart zerod.service`.
Environment=ZEROD_BEARER_TOKEN={{bearer_token}}

# Configuration files (bind address, systemd unit allowlist, managed
# config registry) live in ~/.config/zerod/zerod.toml — see
# `zerod.toml.example` for the schema. Drop additional env overrides
# in ~/.config/systemd/user/zerod.service.d/override.conf to keep them
# safe across `zerod service install --force`.

[Install]
WantedBy=default.target
